# ticketdesk.conf
# ============================================================================
# Supervisor Configuration - TicketDesk Enterprise v2.1
# Ubicar en: /etc/supervisor/conf.d/ticketdesk.conf
# ============================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ============================================================================
# TicketDesk Flask Application
# ============================================================================
[program:ticketdesk]
command=/usr/bin/python3 /opt/ticketdesk/app.py
directory=/opt/ticketdesk
user=ticketdesk
environment=FLASK_ENV=production,PYTHONUNBUFFERED=1,PATH="/usr/local/bin:/usr/bin"

# Restart Policy
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stopwaitsecs=30

# Logging
stdout_logfile=/var/log/ticketdesk/app.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/ticketdesk/app_error.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10

# Process Management
numprocs=1
priority=999

# Healthcheck (restart si proceso muere)
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true

# ============================================================================
# Scheduler / Background Jobs (Opcional)
# ============================================================================
[program:ticketdesk-scheduler]
command=/usr/bin/python3 /opt/ticketdesk/scheduler.py
directory=/opt/ticketdesk
user=ticketdesk
environment=FLASK_ENV=production,PYTHONUNBUFFERED=1

autostart=true
autorestart=true
startsecs=10
stopasgroup=true
stopwaitsecs=30

stdout_logfile=/var/log/ticketdesk/scheduler.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/ticketdesk/scheduler_error.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10

# ============================================================================
# Backup Service (Opcional - Backups diarios)
# ============================================================================
[program:ticketdesk-backup]
command=/bin/bash /opt/ticketdesk/backup.sh
directory=/opt/ticketdesk
user=ticketdesk

autostart=true
autorestart=true
startsecs=10
startretries=3

stdout_logfile=/var/log/ticketdesk/backup.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/ticketdesk/backup_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

# ============================================================================
# Evento: Auto-restart en crash
# ============================================================================
[eventlistener:ticketdesk-watchdog]
command=/usr/bin/python3 -m supervisor.supervisord
buffer_size=100
events=TICK_60,PROCESS_STATE
result_handler=ticketdesk_eventhandler

[eventhandler:ticketdesk_eventhandler]
command=bash -c "echo '%(eventname)s' >> /var/log/ticketdesk/watchdog.log"
buffer_size=100
events=PROCESS_STATE_FATAL,PROCESS_STATE_EXITED_UNEXPECTED
result_handler=ticketdesk_eventhandler

# ============================================================================
# Group - Para iniciar/parar todo junto
# ============================================================================
[group:ticketdesk-services]
programs=ticketdesk,ticketdesk-scheduler,ticketdesk-backup
priority=999
